# services/kimnaty.trend.year.timer
# This timer is for triggering the trending of past year(s) data

[Unit]
Description=daily trending data from past year (timer)

[Timer]
OnCalendar=*-*-* 01:01
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
